# LAB/boot_count/bootcount.service
[Unit]
Description=System Boot count

[Service]
Type=oneshot
ExecStart=/mnt/S8800.sh

[Install]
WantedBy=multi-user.target
